# deployment/systemd/cs-scheduler.service
[Unit]
Description=Chrono Swarm Scheduler
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cs-scheduler
Restart=always
User=chrono

[Install]
WantedBy=multi-user.target
